# core/etc/profile.d/aliases.sh
alias ll="ls -lhac"
alias code="code-server"